# install/90-bt-autoconnect.rules
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/bin/systemctl start bt-autoconnect.service"